# sys-apps/hv_vss_daemon/files/hv_vss_daemon-initd
#!/sbin/runscript

# path to hv_vss_daemon binary
HV_VSS_BIN=/usr/sbin/hv_vss_daemon

depend() {
        use clock logger
        need localmount
        provide hv_vss_daemon
}

start(){
        ebegin "Starting Hyper-V VSS daemon"
        start-stop-daemon --start --quiet --exec $HV_VSS_BIN
        eend $?
}

stop(){
        ebegin "Stopping Hyper-V VSS daemon"
        start-stop-daemon --stop --quiet $HV_VSS_BIN
        eend $?
}

